# systemd/airlock.service
# /etc/systemd/system/airlock.service
#
# THE AIRLOCK v5.0.8 FORTRESS-HARDENED — USB Sanitization Station
#
# mount/umount/deauthorize artık privileged helper üzerinden yapılır.
# Bu servis NON-ROOT çalışır, hiçbir capability gerektirmez.
#
# Kurulum:
#   sudo cp /opt/airlock/systemd/airlock.service /etc/systemd/system/
#   sudo cp /opt/airlock/systemd/airlock-helper.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable airlock-helper.service airlock.service
#   sudo systemctl start airlock-helper.service airlock.service

[Unit]
Description=THE AIRLOCK v5.0.8 FORTRESS-HARDENED — USB Sanitization Station
After=multi-user.target airlock-helper.service
Wants=clamav-daemon.service
Requires=airlock-helper.service

[Service]
Type=simple
User=airlock
Group=airlock
WorkingDirectory=/opt/airlock
ExecStart=/opt/airlock/venv/bin/python3 /opt/airlock/app/main.py
Restart=always
RestartSec=5
WatchdogSec=300
UMask=0077

# ── Ortam Degiskenleri ──
Environment=PYTHONUNBUFFERED=1
Environment=LANG=en_US.UTF-8

# ── Guvenlik Sertlestirme (MAKSIMUM) ──
# Artik CAP_SYS_ADMIN YOK — mount/umount helper uzerinden yapilir
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RemoveIPC=true

# ── Ag Erisimi Engelle (air-gapped) ──
RestrictAddressFamilies=AF_UNIX AF_LOCAL
PrivateNetwork=true

# ── Dosya Sistemi Izinleri ──
ReadWritePaths=/opt/airlock/data /opt/airlock/tmp /opt/airlock/data/logs /opt/airlock/data/quarantine /mnt/airlock_target
ReadOnlyPaths=/opt/airlock/app /opt/airlock/config /opt/airlock/keys /mnt/airlock_source /mnt/airlock_update
TemporaryFileSystem=/opt/airlock/tmp:size=512M,mode=0750

# ── Capabilities: HIC GEREKMIYOR ──
# mount/umount/deauthorize privileged helper servisinde
CapabilityBoundingSet=
AmbientCapabilities=

# ── USB, GPIO, I2C, Ses Erisimi ──
SupplementaryGroups=gpio i2c audio plugdev disk

# ── Kaynak Limitleri ──
LimitNOFILE=65536
LimitNPROC=256
MemoryMax=2G
TasksMax=128

# ── Loglama ──
StandardOutput=journal
StandardError=journal
SyslogIdentifier=airlock

[Install]
WantedBy=multi-user.target
